# scripts/com.casterly.claude-tty-keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.casterly.claude-tty-keepalive</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/casterly/.claude/scripts/claude-tty-keepalive.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>15</integer>
    <key>StandardOutPath</key>
    <string>/tmp/claude-tty-keepalive.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-tty-keepalive.log</string>
</dict>
</plist>
